# docker/php-fpm/supervisord.conf
[program:messenger-consume]
command=php /var/www/bin/console messenger:consume-messages async --limit=100 --time-limit=3600
numprocs=1
autostart=true
autorestart=true
startsecs=10
user=www-data


